# systemd/pi-fm-rds-go.service
[Unit]
Description=Pi FM RDS Go - FM Broadcast Service
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/data/pi-fm-rds
ExecStart=/data/pi-fm-rds/pi-fm-rds-server -config /data/pi-fm-rds/config.yaml
Restart=on-failure
RestartSec=10s

Environment="GIN_MODE=release"

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-fm-rds-go

PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/data/pi-fm-rds /data/pi-fm-rds/uploads /data/pi-fm-rds/cache

LimitNOFILE=65536
LimitNPROC=512

[Install]
WantedBy=multi-user.target
